# script/com.example.applaud.filemover.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.applaud.filemover</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>/Users/m1mac/.local/bin/applaud_launcher.sh</string>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>WatchPaths</key>
    <array>
        <string>/Volumes/NO NAME</string>
    </array>
    <key>StandardOutPath</key>
    <string>/tmp/com.example.applaud.filemover.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.example.applaud.filemover.stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/bin:/usr/bin:/usr/local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/Users/m1mac/Library/Python/3.9/bin</string>
    </dict>
    <key>SessionCreate</key>
    <true/>
</dict>
</plist>
